# supervisord.conf
[supervisord]
logfile=log/supervisord.log ; main log file; default $CWD/supervisord.log
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:gunicorn-django]
command=gunicorn izone.wsgi -b 0.0.0.0:8000
directory=/opt/cloud/izone
stdout_logfile=log/gunicorn.django.log
stderr_logfile=log/gunicorn.django.log
autostart=true
autorestart=true
startsecs=2
stopwaitsecs=2
priority=100

[program:celery-worker]
command=celery -A izone worker -l info
directory=/opt/cloud/izone
stdout_logfile=log/celery.worker.log
stderr_logfile=log/celery.worker.log
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=3
priority=102

[program:celery-beat]
command=celery -A izone beat -l info
directory=/opt/cloud/izone
stdout_logfile=log/celery.beat.log
stderr_logfile=log/celery.beat.log
autostart=true
autorestart=true
startsecs=2
stopwaitsecs=2
priority=101